# emuleted_hue_admin/rootfs/etc/s6-overlay/s6-rc.d/emulated-hue-manager/finish
#!/command/with-contenv bashio
# ==============================================================================
# Called when uvicorn exits
# ==============================================================================
bashio::log.warning "Emulated Hue Manager exited with code ${1}"
