# emuleted_hue_admin/rootfs/etc/s6-overlay/s6-rc.d/emulated-hue-manager/run
#!/command/with-contenv bashio
# ==============================================================================
# Start the Emulated Hue Manager (uvicorn)
# ==============================================================================
bashio::log.info "Starting Emulated Hue Manager..."

export PYTHONPATH=/

exec uvicorn app.main:app \
    --host 0.0.0.0 \
    --port 8099 \
    --workers 1 \
    --log-level info
